# com.macrobrief.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.macrobrief.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/brendanberthold/Dropbox/My Computer/GitHub/TheGreatInformationGatherer/.venv/bin/python3</string>
        <string>/Users/brendanberthold/Dropbox/My Computer/GitHub/TheGreatInformationGatherer/run_daily_pipeline.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/brendanberthold/Dropbox/My Computer/GitHub/TheGreatInformationGatherer</string>

    <key>StartCalendarInterval</key>
    <array>
        <!-- Run at 14:00 UTC (9:00 AM ET / 3:00 PM CET) Mon-Fri -->
        <dict>
            <key>Weekday</key><integer>1</integer>
            <key>Hour</key><integer>14</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>2</integer>
            <key>Hour</key><integer>14</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>3</integer>
            <key>Hour</key><integer>14</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>4</integer>
            <key>Hour</key><integer>14</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>5</integer>
            <key>Hour</key><integer>14</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/brendanberthold/Dropbox/My Computer/GitHub/TheGreatInformationGatherer/logs/daily_pipeline.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/brendanberthold/Dropbox/My Computer/GitHub/TheGreatInformationGatherer/logs/daily_pipeline_error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
